# livekit-examples/cloud-gaming/steam-headless/overlay/etc/supervisor.d/gstreamer.ini
[program:gstreamer]
priority=30
autostart=true
autorestart=true
directory=/
user=root
command=gst-launch-1.0 ximagesrc ! video/x-raw ! videoscale ! video/x-raw,width=1280,height=720 ! videoconvert ! capsfilter caps="video/x-raw,format=I420" ! nvh264enc bitrate=1500 ! h264parse config-interval=-1 ! queue max-size-buffers=1 ! livekit.video_0 pulsesrc ! audioconvert ! opusenc ! queue ! livekit.audio_0 livekitwebrtcsink name=livekit signaller::ws-url=%(ENV_LK_ENDPOINT)s signaller::room-name=%(ENV_LK_ROOM)s signaller::api-key=%(ENV_LK_API_KEY)s signaller::secret-key=%(ENV_LK_SECRET_KEY)s
stopsignal=INT
stdout_logfile=/var/log/supervisor/gstreamer.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/log/supervisor/gstreamer.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
environment=DISPLAY=":55"

