# livekit-examples/cloud-gaming/steam-headless/overlay/etc/supervisor.d/input-handler.ini
[program:input-handler]
priority=2
autostart=true
autorestart=true
directory=/usr/bin
user=root
command= input-handler
stopsignal=INT
stdout_logfile=/var/log/supervisor/input-handler.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/log/supervisor/input-handler.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
